# deploy/systemd/labuse.service
# =============================================================================
# LA BUSE — service systemd (FastAPI/Uvicorn derrière Nginx)
# -----------------------------------------------------------------------------
# Pose (NE PAS appliquer maintenant — voir docs/DEPLOYMENT_OVH_VPS.md §8) :
#   install -o root -g root -m 644 deploy/systemd/labuse.service /etc/systemd/system/labuse.service
#   systemctl daemon-reload && systemctl enable --now labuse
#   journalctl -u labuse -f          # logs (journald — JAMAIS de nohup)
# =============================================================================
[Unit]
Description=LA BUSE — radar foncier (Uvicorn)
Documentation=file:/opt/labuse/app/docs/DEPLOYMENT_OVH_VPS.md
After=network-online.target postgresql.service
Wants=network-online.target
Requires=postgresql.service

[Service]
Type=exec
User=labuse
Group=labuse
WorkingDirectory=/opt/labuse/app

# Secrets & configuration (LABUSE_* + ANTHROPIC_API_KEY) — hors du dépôt git.
EnvironmentFile=/etc/labuse/labuse.env

# Uvicorn, 2 workers (équilibre 4 vCores partagés avec PostgreSQL ; montable à 3 si besoin).
# --proxy-headers : on est derrière Nginx ; --forwarded-allow-ips limite la confiance au proxy local.
ExecStart=/opt/labuse/venv/bin/uvicorn labuse.api.app:app \
    --host 127.0.0.1 --port 8000 \
    --workers 2 \
    --proxy-headers --forwarded-allow-ips 127.0.0.1

# Redémarrage automatique en cas de crash (jamais de service mort silencieux).
Restart=on-failure
RestartSec=3
TimeoutStartSec=60

# Logs → journald.
StandardOutput=journal
StandardError=journal
SyslogIdentifier=labuse

# Durcissement raisonnable (l'app n'écrit rien hors /opt/labuse et /tmp privé).
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictSUIDSGID=true

[Install]
WantedBy=multi-user.target
