# deploy/cron.d/sitadel
# LA BUSE — SITADEL quotidien-en-attente-du-mensuel (J+2) : le refresh delta (recouvrement
# 3 mois) est un no-op léger tant que le SDES n'a pas publié ; le jour de la livraison
# mensuelle, l'ingestion part le jour même. Puis dérivés (hebdo le lundi : + m10 vélocité).
SHELL=/bin/sh
15 4 * * * labuse set -a; . /etc/labuse/labuse.env; set +a; cd /opt/labuse/app && /opt/labuse/venv/bin/python -m labuse.ingestion.permits_sdes --refresh >> /var/log/labuse/sitadel_refresh.log 2>&1 && /opt/labuse/venv/bin/labuse fraicheur-derives >> /var/log/labuse/sitadel_refresh.log 2>&1
45 5 * * 1 labuse set -a; . /etc/labuse/labuse.env; set +a; cd /opt/labuse/app && /opt/labuse/venv/bin/labuse fraicheur-derives --hebdo >> /var/log/labuse/sitadel_refresh.log 2>&1
