# deploy/cron.d/solaire
# Habitat Solaire — refresh MENSUEL des sources vivantes du module (le gisement PVGIS,
# lui, est climatologique : one-shot, pas de refresh) :
#   - registre national des installations 974 (repowering, densité communale)
#   - capacités d'accueil réseau EDF SEI (MAJ trimestrielle publiée)
#   - parkings APER (recalcul tranches/échéances + signaux deadline)
#   - flags qualification + conso/facture (DPE et DVF bougent)
#   - vue tertiaire (dépend de tout ce qui précède)
#   - purge STRICTE du cache Google Solar API > 30 j (conformité ToS, Lot 8)
# Passage le 6 à 04h00 — après la BAN (le 5 à 03h30).
#
# Installation VPS :  install -o root -g root -m 644 deploy/cron.d/solaire /etc/cron.d/solaire
SHELL=/bin/sh
0 4 6 * * labuse set -a; . /etc/labuse/labuse.env; set +a; cd /opt/labuse/app && /opt/labuse/venv/bin/labuse solaire-pv-registry >> /var/log/labuse/solaire_refresh.log 2>&1 && /opt/labuse/venv/bin/labuse solaire-grid-capacity >> /var/log/labuse/solaire_refresh.log 2>&1 && /opt/labuse/venv/bin/labuse solaire-parkings >> /var/log/labuse/solaire_refresh.log 2>&1 && /opt/labuse/venv/bin/labuse solaire-flags >> /var/log/labuse/solaire_refresh.log 2>&1 && /opt/labuse/venv/bin/labuse solaire-conso >> /var/log/labuse/solaire_refresh.log 2>&1 && /opt/labuse/venv/bin/labuse solaire-tertiaire >> /var/log/labuse/solaire_refresh.log 2>&1 && /opt/labuse/venv/bin/labuse solaire-cache-purge >> /var/log/labuse/solaire_refresh.log 2>&1
